# input/deploy/macos/com.myhron.eastmoney-yw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.myhron.eastmoney-yw</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/moyueheng/.openclaw/workspace/input/scripts/run-eastmoney-yw.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/moyueheng/.openclaw/workspace/input</string>
    <key>StandardOutPath</key>
    <string>/Users/moyueheng/.openclaw/workspace/input/data/state/eastmoney-yw.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/moyueheng/.openclaw/workspace/input/data/state/eastmoney-yw.stderr.log</string>
  </dict>
</plist>
